# server/supervisor/scheduler.conf
[program:scheduler]
command=/app/uptimatic scheduler
user=uptimatic
autostart=true
autorestart=true
stdout_logfile=/var/log/scheduler.log
stderr_logfile=/var/log/scheduler.err
