# server/supervisor/server.conf
[program:server]
command=/app/uptimatic server
user=uptimatic
autostart=true
autorestart=true
stdout_logfile=/var/log/server.log
stderr_logfile=/var/log/server.err
