# server/supervisor/worker.conf
[program:worker]
command=/app/uptimatic worker
autostart=true
autorestart=true
stdout_logfile=/var/log/worker.log
stderr_logfile=/var/log/worker.err
