# server/s6/scheduler/run
#!/bin/sh
exec /app/uptimatic scheduler
